# s6/services.d/rudolfs/run
#!/usr/bin/with-contenv sh

set -eou pipefail

fdmove -c 2 1 \
  rudolfs \
  --host=${RUDOLFS_ADDR} \
  --key=${RUDOLFS_KEY} \
  local \
  --path=${RUDOLFS_PATH}
