# oneliner/etc/s6-overlay/s6-rc.d/rudolfs/run
#!/usr/bin/env bash

set -euo pipefail

# Port 8888 - Git LFS Server

exec 2>&1
exec rudolfs \
	--host=127.0.0.1:8888 \
	--key=$(cat /var/data/rudolfs/key) \
	local \
	--path=/var/data/rudolfs/data
