# oneliner/etc/s6-overlay/s6-rc.d/api-log/run
#!/bin/bash

set -eou pipefail

sed 's/^/[api] /'
