# oneliner/etc/s6-overlay/s6-rc.d/api-log/run
#!/usr/bin/env bash

set -euo pipefail

sed 's/^/[api] /'
