# oneliner/etc/s6-overlay/s6-rc.d/postgresql/run
#!/command/execlineb -P

fdmove -c 2 1
s6-setuidgid postgres
/usr/lib/postgresql/14/bin/postgres -D "/var/data/postgresql/data"
