# oneliner/etc/s6-overlay/s6-rc.d/rudolfs-prepare/up
/bin/bash /etc/scripts/rudolfs/setup.sh
