# s6/services.d/backend/finish
#!/usr/bin/execlineb -S1

# if backend stops, it brings the whole container down
s6-svscanctl -t /var/run/s6/services
