# oneliner/etc/s6-overlay/s6-rc.d/postgresql-log/run
#!/bin/bash

set -eou pipefail

sed 's/^/[postgresql] /'
