# oneliner/etc/s6-overlay/s6-rc.d/rudolfs-log/run
#!/usr/bin/env bash

set -euo pipefail

sed --unbuffered 's/^/[rudolfs] /'
