# oneliner/etc/s6-overlay/s6-rc.d/rudolfs-log/run
#!/bin/bash

set -eou pipefail

sed 's/^/[rudolfs] /'
