# oneliner/etc/s6-overlay/s6-rc.d/rudolfs-log/run
#!/usr/bin/env sh

set -eou pipefail

sed 's/^/[rudolfs] /'
