# oneliner/etc/s6-overlay/s6-rc.d/reproxy/run
#!/usr/bin/env sh

set -eou pipefail

exec 2>&1
exec reproxy \
  --assets.location='/web/dist' \
  --assets.spa \
  --static.enabled \
  --static.rule='*,^/api/(.*),http://127.0.0.1:3000/$1,' \
  --gzip \
  --logger.stdout \
  --listen="0.0.0.0:80"
