# oneliner/etc/s6-overlay/s6-rc.d/postgresql-prepare/up
/bin/sh /etc/scripts/postgresql/setup.sh
