# oneliner/etc/s6-overlay/s6-rc.d/reproxy-log/run
#!/bin/bash

set -eou pipefail

sed 's/^/[reproxy] /'
